# configs/shared-mcp/shared-mcp.service
[Unit]
Description=Shared MCP Hub (Podman Compose)
Documentation=https://github.com/dinesh-nimmagadda3/my-ai-tools
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=$REPLACE_WITH_PROJECT_DIR/configs/shared-mcp
Environment=PROJECT_DIR=$REPLACE_WITH_PROJECT_DIR
# The compose file uses host volumes, so we must expose these variables
ExecStart=/usr/bin/podman-compose up -d
ExecStop=/usr/bin/podman-compose down
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
