# configs/shared-mcp/shared-mcp-pod.service
[Unit]
Description=Shared MCP Hub Pod (Podman Kube Play)
Documentation=https://github.com/dinesh-nimmagadda3/my-ai-tools
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# We use absolute paths substituted by cli.sh to avoid systemd variable expansion issues
ExecStartPre=-/usr/bin/podman kube down $REPLACE_WITH_PROJECT_DIR/configs/shared-mcp/hub-pod.yaml
# We rely on the hostPort definition inside hub-pod.yaml instead of --publish
ExecStart=/usr/bin/podman kube play --network bridge $REPLACE_WITH_PROJECT_DIR/configs/shared-mcp/hub-pod.yaml
ExecStop=/usr/bin/podman kube down $REPLACE_WITH_PROJECT_DIR/configs/shared-mcp/hub-pod.yaml
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
