# configs/shared-mcp/shared-mcp.service
[Unit]
Description=Shared MCP Hub (Podman Compose)
Documentation=https://github.com/dinesh-nimmagadda3/my-ai-tools
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=$REPLACE_WITH_PROJECT_DIR/configs/shared-mcp
Environment=PROJECT_DIR=$REPLACE_WITH_PROJECT_DIR
# Use blocking 'up' so systemd can track the process and logs
ExecStart=/usr/bin/podman compose up
ExecStop=/usr/bin/podman compose down
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
